# deploy/vyrdx-relay.service
# VYRDX Relay — systemd service for DELL side
# Exposes /opt/vyrdx state files over HTTP on the Tailscale mesh.
# Runs as root to read vyrdon-owned state files.
#
# Install:
#   sudo cp deploy/vyrdx-relay.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vyrdx-relay

[Unit]
Description=VYRDX Relay — HTTP bridge for /opt/vyrdx state
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
WorkingDirectory=/opt/kitty
ExecStart=/usr/bin/npx tsx server/vyrdx-relay.ts
Restart=always
RestartSec=5

Environment=VYRDX_CORE_BASE=/opt/vyrdx/core
Environment=VYRDX_RELAY_PORT=7801
Environment=VYRDX_RELAY_HOST=0.0.0.0
Environment=NODE_ENV=production

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vyrdx-relay

[Install]
WantedBy=multi-user.target
