# deploy/vyrdx-tunnel.service
# SSH reverse tunnel: forwards droplet localhost:7801 → ASUS localhost:7801
# So VXSTATION on the droplet can reach the VYRDX relay on ASUS.
#
# Install:
#   sudo cp deploy/vyrdx-tunnel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vyrdx-tunnel

[Unit]
Description=VYRDX Relay SSH Tunnel — ASUS to Droplet
After=network-online.target vyrdx-relay.service
Wants=network-online.target vyrdx-relay.service

[Service]
Type=simple
User=t79
ExecStart=/usr/bin/ssh -o IdentitiesOnly=yes -i /home/t79/.ssh/id_deploy_vyrdon \
  -o ServerAliveInterval=30 -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes -o StrictHostKeyChecking=no \
  -R 7801:127.0.0.1:7801 -N root@134.199.227.138
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
